# scripts/nginx.conf
[program:nginx]
command=/usr/sbin/nginx
autostart=true
autoretry=true
startretries=15
